# UNIX/DD_CONF/supervisor/bull2017.conf
[group:bull2017alloc]
programs=bull2017alloc-1,bull2017alloc-2,bull2017alloc-3,bull2017alloc-4

[program:bull2017alloc-1]
command=/home/Output/BullFight2017/AllocServer/bin/BullFight2017Alloc -p 23401 -s 23401 -l 1
directory=/home/Output/BullFight2017/AllocServer/bin
user=root
priority=700

[program:bull2017alloc-2]
command=/home/Output/BullFight2017/AllocServer/bin/BullFight2017Alloc -p 23402 -s 23402 -l 2
directory=/home/Output/BullFight2017/AllocServer/bin
user=root
priority=700

[program:bull2017alloc-3]
command=/home/Output/BullFight2017/AllocServer/bin/BullFight2017Alloc -p 23403 -s 23403 -l 3
directory=/home/Output/BullFight2017/AllocServer/bin
user=root
priority=700

[program:bull2017alloc-4]
command=/home/Output/BullFight2017/AllocServer/bin/BullFight2017Alloc -p 23404 -s 23404 -l 4
directory=/home/Output/BullFight2017/AllocServer/bin
user=root
priority=700




[group:bull2017game]
programs=bull2017game-1,bull2017game-2,bull2017game-3,bull2017game-4

[program:bull2017game-1]
command=/home/Output/BullFight2017/GameServer/bin/BullFight2017Game -p 23501 -s 23501 -l 1
directory=/home/Output/BullFight2017/GameServer/bin
user=root
priority=800

[program:bull2017game-2]
command=/home/Output/BullFight2017/GameServer/bin/BullFight2017Game -p 23502 -s 23502 -l 2
directory=/home/Output/BullFight2017/GameServer/bin
user=root
priority=800

[program:bull2017game-3]
command=/home/Output/BullFight2017/GameServer/bin/BullFight2017Game -p 23503 -s 23503 -l 3
directory=/home/Output/BullFight2017/GameServer/bin
user=root
priority=800

[program:bull2017game-4]
command=/home/Output/BullFight2017/GameServer/bin/BullFight2017Game -p 23504 -s 23504 -l 4
directory=/home/Output/BullFight2017/GameServer/bin
user=root
priority=800




[group:bull2017robot]
programs=bull2017robot-11,bull2017robot-12,bull2017robot-21,bull2017robot-22,bull2017robot-31,bull2017robot-32,bull2017robot-41,bull2017robot-42



[program:bull2017robot-11]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23601 -s 23601 -l 1 -n 1
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900

[program:bull2017robot-12]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23602 -s 23602 -l 2 -n 2
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900



[program:bull2017robot-21]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23603 -s 23603 -l 1 -n 3
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900

[program:bull2017robot-22]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23604 -s 23604 -l 2 -n 4
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900



[program:bull2017robot-31]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23605 -s 23605 -l 1 -n 5
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900

[program:bull2017robot-32]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23606 -s 23606 -l 2 -n 6
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900



[program:bull2017robot-41]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23607 -s 23607 -l 1 -n 7
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900

[program:bull2017robot-42]
command=/home/Output/BullFight2017/RobotServer/bin/BullFight2017Robot -p 23608 -s 23608 -l 2 -n 8
directory=/home/Output/BullFight2017/RobotServer/bin
user=root
priority=900