# UNIX/DD_CONF/supervisor/flower2017.conf
[group:flower2017alloc]
programs=flower2017alloc-1,flower2017alloc-2,flower2017alloc-3,flower2017alloc-4

[program:flower2017alloc-1]
command=/home/Output/FlowerGame2017/AllocServer/bin/Flower2017Alloc -p 28101 -s 28101 -l 1
directory=/home/Output/FlowerGame2017/AllocServer/bin
user=root
priority=700

[program:flower2017alloc-2]
command=/home/Output/FlowerGame2017/AllocServer/bin/Flower2017Alloc -p 28102 -s 28102 -l 2
directory=/home/Output/FlowerGame2017/AllocServer/bin
user=root
priority=700

[program:flower2017alloc-3]
command=/home/Output/FlowerGame2017/AllocServer/bin/Flower2017Alloc -p 28103 -s 28103 -l 3
directory=/home/Output/FlowerGame2017/AllocServer/bin
user=root
priority=700

[program:flower2017alloc-4]
command=/home/Output/FlowerGame2017/AllocServer/bin/Flower2017Alloc -p 28104 -s 28104 -l 4
directory=/home/Output/FlowerGame2017/AllocServer/bin
user=root
priority=700


[group:flower2017game]
programs=flower2017game-1,flower2017game-2,flower2017game-3,flower2017game-4

[program:flower2017game-1]
command=/home/Output/FlowerGame2017/GameServer/bin/Flower2017Game -p 28201 -l 1 -s 28201
directory=/home/Output/FlowerGame2017/GameServer/bin
user=root
priority=800

[program:flower2017game-2]
command=/home/Output/FlowerGame2017/GameServer/bin/Flower2017Game -p 28202 -l 2 -s 28202
directory=/home/Output/FlowerGame2017/GameServer/bin
user=root
priority=800

[program:flower2017game-3]
command=/home/Output/FlowerGame2017/GameServer/bin/Flower2017Game -p 28203 -l 3 -s 28203
directory=/home/Output/FlowerGame2017/GameServer/bin
user=root
priority=800

[program:flower2017game-4]
command=/home/Output/FlowerGame2017/GameServer/bin/Flower2017Game -p 28204 -l 4 -s 28204
directory=/home/Output/FlowerGame2017/GameServer/bin
user=root
priority=800


[group:flower2017robot]
programs=flower2017robot-1,flower2017robot-2

[program:flower2017robot-1]
command=/home/Output/FlowerGame2017/RobotServer/bin/Flower2017Robot -p 28301 -s 28301 -l 1 -n 1
directory=/home/Output/FlowerGame2017/RobotServer/bin
user=root
priority=900

[program:flower2017robot-2]
command=/home/Output/FlowerGame2017/RobotServer/bin/Flower2017Robot -p 28302 -s 28302 -l 1 -n 2
directory=/home/Output/FlowerGame2017/RobotServer/bin
user=root
priority=900
