# UNIX/DD_CONF/supervisor/land2017.conf
[group:land2017alloc]
programs=Land2017Alloc-1,Land2017Alloc-2,Land2017Alloc-3,Land2017Alloc-4


[program:Land2017Alloc-1]
command=/home/Output/LandGame2017/AllocServer/bin/Land2017Alloc -s 23701 -p 23701 -l 1
directory=/home/Output/LandGame2017/AllocServer/bin
user=root
priority=700

[program:Land2017Alloc-2]
command=/home/Output/LandGame2017/AllocServer/bin/Land2017Alloc -s 23702 -p 23702 -l 2
directory=/home/Output/LandGame2017/AllocServer/bin
user=root
priority=700

[program:Land2017Alloc-3]
command=/home/Output/LandGame2017/AllocServer/bin/Land2017Alloc -s 23703 -p 23703 -l 3
directory=/home/Output/LandGame2017/AllocServer/bin
user=root
priority=700

[program:Land2017Alloc-4]
command=/home/Output/LandGame2017/AllocServer/bin/Land2017Alloc -s 23704 -p 23704 -l 4
directory=/home/Output/LandGame2017/AllocServer/bin
user=root
priority=700


[group:land2017game]
programs=land2017game-1,land2017game-2,land2017game-3,land2017game-4

[program:land2017game-1]
command=/home/Output/LandGame2017/GameServer/bin/Land2017Game -s 23801 -p 23801 -l 1
directory=/home/Output/LandGame2017/GameServer/bin
user=root
priority=800

[program:land2017game-2]
command=/home/Output/LandGame2017/GameServer/bin/Land2017Game -s 23802 -p 23802 -l 2
directory=/home/Output/LandGame2017/GameServer/bin
user=root
priority=800

[program:land2017game-3]
command=/home/Output/LandGame2017/GameServer/bin/Land2017Game -s 23803 -p 23803 -l 3
directory=/home/Output/LandGame2017/GameServer/bin
user=root
priority=800

[program:land2017game-4]
command=/home/Output/LandGame2017/GameServer/bin/Land2017Game -s 23804 -p 23804 -l 4
directory=/home/Output/LandGame2017/GameServer/bin
user=root
priority=800


[group:land2017robot]
programs=land2017robot-1,land2017robot-2,land2017robot-3,land2017robot-4

[program:land2017robot-1]
command=/home/Output/LandGame2017/RobotServer/bin/Land2017Robot -s 23901 -p 23901 -l 1
directory=/home/Output/LandGame2017/RobotServer/bin
user=root
priority=900

[program:land2017robot-2]
command=/home/Output/LandGame2017/RobotServer/bin/Land2017Robot -s 23902 -p 23902 -l 2
directory=/home/Output/LandGame2017/RobotServer/bin
user=root
priority=900

[program:land2017robot-3]
command=/home/Output/LandGame2017/RobotServer/bin/Land2017Robot -s 23903 -p 23903 -l 3
directory=/home/Output/LandGame2017/RobotServer/bin
user=root
priority=900

[program:land2017robot-4]
command=/home/Output/LandGame2017/RobotServer/bin/Land2017Robot -s 23904 -p 23904 -l 4
directory=/home/Output/LandGame2017/RobotServer/bin
user=root
priority=900


