# UNIX/config/supervisor/httpsvr.conf
[program:httpsvr]
command=/home/Output/Http/bin/HttpSvr
directory=/home/Output/Http/bin
user=root

